# installer/systemd/sky-daemon.service
[Unit]
Description=Sky Panel daemon (sky-daemon)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
EnvironmentFile=/opt/sky-panel/sky-daemon.env
ExecStart=/opt/sky-panel/bin/sky-daemon
Restart=on-failure
RestartSec=2
# The daemon talks to the Docker Engine API over /var/run/docker.sock, so it
# needs to run as a user in the docker group rather than fully unprivileged.
User=root

[Install]
WantedBy=multi-user.target
